# docker/web/supervisord.conf
[supervisord]
nodaemon=true

[program:mysqld]
command=/etc/init.d/mariadb start

[program:activemq]
# Start the standalone broker only when the web app is NOT hosting an embedded one.
# With JMS_BROKER_EMBEDDED=true (the image default) the web process runs the broker
# itself, so this block exits immediately and supervisord leaves it stopped — saving
# a whole broker JVM (~150 MB+). Set JMS_BROKER_EMBEDDED=false to use this instead.
command=bash -c 'if [ "${JMS_BROKER_EMBEDDED:-true}" = "true" ]; then echo "Embedded broker in use; standalone ActiveMQ skipped."; exit 0; fi; exec /opt/activemq/bin/activemq start'
# Cap the broker heap (ActiveMQ defaults to -Xmx1g) for the small-demo footprint.
environment=ACTIVEMQ_OPTS_MEMORY="-Xms32m -Xmx128m"
# The conditional command above exits 0 when embedded; treat that as success, not a crash.
startsecs=0
autorestart=false

[program:voj-web]
# Lean JVM defaults for the small-demo footprint: capped heap/metaspace, the
# single-threaded GC and no C2 JIT. Keep this in sync with scripts/run-web.sh.
command=bash -c 'sleep 15 && java -Xmx256m -XX:+UseSerialGC -XX:MaxMetaspaceSize=128m -XX:TieredStopAtLevel=1 -Xss512k -jar /opt/voj/voj.web.jar'
directory=/opt/voj

